# deploy/etcd/theia-etcd.service
[Unit]
Description=Theia etcd — single-node state backbone
Documentation=docs/tasks/BACKLOG/etcd-state-backbone.md
After=network.target docker.service
Wants=network-online.target

[Service]
Type=notify
User=etcd
Group=etcd
# Listen on two interfaces so both the host (supdbg, etcdctl, GUI on
# host) and the docker compose containers (supervisor, services/com,
# eventually services/db) share one etcd instance:
#
#   127.0.0.1:2379    workspace host
#   172.30.0.1:2379   theia_net docker bridge gateway IP
#
# If you redefine the bridge subnet in deploy/docker-compose.yml,
# update both this URL and the README. Containers connect with
# `--endpoints=172.30.0.1:2379`.
ExecStart=/usr/bin/etcd \
    --name theia0 \
    --data-dir /var/lib/theia-etcd \
    --listen-client-urls http://127.0.0.1:2379,http://172.30.0.1:2379 \
    --advertise-client-urls http://127.0.0.1:2379 \
    --listen-peer-urls http://127.0.0.1:2380 \
    --initial-advertise-peer-urls http://127.0.0.1:2380 \
    --initial-cluster theia0=http://127.0.0.1:2380 \
    --initial-cluster-state new \
    --initial-cluster-token theia-etcd-token \
    --log-output stderr \
    --quota-backend-bytes 268435456
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
