# deploy/theia-supervisor.service
[Unit]
Description=Theia supervisor (boots the supervised service tree)
Documentation=https://github.com/perotheia/theia
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# theia-run exports the strict env (THEIA_ROOT_DIR + THEIA_SUPERVISOR_MANIFEST)
# the supervisor requires, then exec's it. THEIA_ROOT defaults to /opt/theia;
# override here for a non-default prefix.
Environment=THEIA_ROOT=/opt/theia
ExecStartPre=/sbin/modprobe tipc
ExecStart=/opt/theia/bin/theia-run
# The supervisor owns its children (PR_SET_PDEATHSIG); a clean stop SIGTERMs it
# and it tears the tree down. Give the gRPC/etcd children room to exit.
KillMode=mixed
TimeoutStopSec=15
Restart=on-failure
RestartSec=3
# Run as root: the supervisor fork/execs children + some FCs need caps
# (nm: netlink; nftables/bpf). Tighten per-service later if desired.
User=root

[Install]
WantedBy=multi-user.target
